# com.mikeybee.cortex2.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mikeybee.cortex2</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/bard/Code/cortex_2/launch_api.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/bard/Code/cortex_2</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>CORTEX_CONFIG</key>
        <string>/Users/bard/Code/cortex_2/config/cortex.yaml</string>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/bard/Code/cortex_2/logs/cortex.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/bard/Code/cortex_2/logs/cortex-error.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
